# docker/conf/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/var/supervisord.pid

[program:go]
command =
user = root
autostart = true